# dist/discord-notify.service
[Unit]
Description=Forwards notifications from Discord through Telegram.
After=network-online.target
Requires=network-online.target

[Service]
ExecStart=/usr/bin/discord-notify
Restart=always

[Install]
WantedBy=multi-user.target
